# etc/systemd/system/flask-app.service
[Unit]
Description=Flask App
After=network.target

[Service]
User=zero
WorkingDirectory=/home/zero/develop/raspberrypi-zero/flask-app
ExecStart=/home/zero/develop/raspberrypi-zero/flask-app/.venv/bin/python app.py
Restart=always

[Install]
WantedBy=multi-user.target

